# services/linux/rlsd-server.service
[Unit]
Description=RetroLookingStatisticsDisplay
After=network.target

[Service]
ExecStart=/usr/bin/rlsd -c
ExecStop=/usr/bin/pkill rlsd

[Install]
WantedBy=default.target